# image-files/etc/rc.d/S48flashradio
#!/bin/sh /etc/rc.common

START=48

start() {
        cd /mesh-extender
	/mesh-extender/flash-radio
}

stop() {
        echo
}
