# image-files/etc/rc.d/S48flashradio
#!/bin/sh /etc/rc.common

START=48

start() {
        cd /mesh-extender
	/mesh-extender/flash900 /mesh-extender/radio~rfd900a.ihx /dev/ttyATH0 >/dev/null 2>&1
}

stop() {
        echo
}
