# image-files/etc/rc.d/S48flashradio
#!/bin/sh /etc/rc.common

START=48

start() {
	/mesh-extender/flash-radio
}

stop() {
}
